# _doc/thumbor/supervisord-thumbor.ini
[program:thumbor]
command=thumbor --port=1000%(process_num)s
process_name=thumbor1000%(process_num)s
directory=/home/thumbor/
numprocs=2
user=nginx
autostart=true
autorestart=true
startretries=3
stdout_logfile=/home/_logs/supervisor/thumbor1000%(process_num)s.stdout.log
stderr_logfile=/home/_logs/supervisor/thumbor1000%(process_num)s.stderr.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
